# dotfiles/crontab/.cron
@reboot /home/hydra/Documents/nixdots/dotfiles/crontab/sync-ipaddress.sh
0 * * * * /home/hydra/Documents/nixdots/dotfiles/crontab/sync-ipaddress.sh
